# Tek2/Octopus/roles/result/files/result.service
[Unit]
Description=Result web application
After=network.target

[Service]
User=ubuntu
ExecStart=/usr/bin/npm start
WorkingDirectory=/opt/result
Restart=always

[Install]
WantedBy=multi-user.target
